# project/config/production/crontab/layui_mvc_frame
* * * * * root /bin/bash /var/www/layui_mvc_frame/project/tool/production/check_update.sh
